# extra/first_boot.service
[Unit]
Description=Configure system on first boot.

[Service]
ExecStart=/usr/local/simple-cdd/first_boot.sh

[Install]
WantedBy=multi-user.target
